# perhost/closettop.crontab
0 0 * * * cd /home/argonightly/argouml-gen && perhost/closettop.sh
